# packages/sthp/debian/preinst
#!/bin/bash
set -e

if ! getent group sthp >/dev/null 2>&1; then
    addgroup --system --quiet sthp
fi

if ! getent passwd sthp >/dev/null 2>&1; then
    adduser --system --quiet \
        --home /nonexistent \
        --no-create-home \
        --shell /usr/sbin/nologin \
        --ingroup sthp \
        --disabled-password \
        --gecos "STHP service user" \
        sthp
fi

if [ -f /etc/systemd/system/sthp.service ]; then
    systemctl stop sthp.service >/dev/null 2>&1 || true
    systemctl disable sthp.service >/dev/null 2>&1 || true
fi

exit 0
